# deploy/lead-api.service.example
[Unit]
Description=Lead API
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/sites/site/server
EnvironmentFile=/var/www/sites/site/server/.env
ExecStart=/usr/bin/node /var/www/sites/site/server/src/index.js
Restart=on-failure
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
