# deploy/odi-leads.service
[Unit]
Description=ODI Lead API
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/sites/odi/server
EnvironmentFile=/var/www/sites/odi/server/.env
ExecStart=/usr/bin/node /var/www/sites/odi/server/src/index.js
Restart=on-failure
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
